# crontab
* * * * * python3 /root/task.py >> /root/service.log 2>> /root/service.log
